# gpuworker/agent/config/visionsuit-gpu-agent.service
[Unit]
Description=VisionSuit GPU Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=visionsuit
Group=visionsuit
WorkingDirectory=/opt/visionsuit-gpu-agent
Environment="VISION_SUITE_AGENT_CONFIG=/etc/visionsuit-gpu-agent/config.yaml"
ExecStart=/opt/visionsuit-gpu-agent/venv/bin/uvicorn main:app --host 0.0.0.0 --port 8081
Restart=on-failure
RestartSec=5
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
